# pdf_streaming.socket
[Unit]
Description=pdf streaming api socket

[Socket]
ListenStream=/run/pdf-streaming.sock

[Install]
WantedBy=sockets.target